# .config/autostart/xrdb.desktop
[Desktop Entry]
Type=Application
Name=X resources loader
Exec=xrdb -merge .Xresources
Comment=Load X resource files on startup of DE
